# docker/s6-overlay/s6-rc.d/queue/run
#!/command/execlineb -P
with-contenv
/usr/local/bin/php /var/www/html/artisan queue:work
